# onshape-extension/legacy/onshape-cadgen.service
[Unit]
Description=3dprint-pipeline Onshape Extension Backend
After=network.target

[Service]
Type=simple
User=tommaso
WorkingDirectory=/home/tommaso/projects/3dprint-pipeline/onshape-extension
ExecStart=/home/tommaso/.local/bin/uvicorn backend.app:app --host 127.0.0.1 --port 8420
Restart=on-failure
RestartSec=5
Environment=HOME=/home/tommaso
Environment=PATH=/home/tommaso/.local/bin:/usr/bin:/bin
Environment=CLAUDE_MODEL=sonnet
Environment=PORT=8420

[Install]
WantedBy=multi-user.target
